# onedriver@.service
[Unit]
Description=Mounts onedrive at a target directory.

[Service]
ExecStart=/usr/bin/onedriver %h/%I

[Install]
WantedBy=default.target
