# onedriver@.service
[Unit]
Description=onedriver

[Service]
ExecStart=/usr/bin/onedriver %I
ExecStop=/usr/bin/fusermount -uz %I
Restart=always
RestartSec=1s

[Install]
WantedBy=default.target
