# onedriver@.service
[Unit]
Description=onedriver

[Service]
ExecStart=/usr/bin/onedriver %h/%I

[Install]
WantedBy=default.target
