# onedriver@.service
[Unit]
Description=onedriver

[Service]
# there is a race condition on session start as the computer network goes from offline to
# online that messes up the initial mount if done immediately on login
ExecStartPre=/usr/bin/sleep 10
ExecStart=/usr/bin/onedriver %I
ExecStop=/usr/bin/fusermount -uz %I
Restart=always

[Install]
WantedBy=default.target
