# data/onedriver@.service
[Unit]
Description=Mounts onedrive at a target directory.

[Service]
ExecStart=onedriver %h/%I

[Install]
WantedBy=default.target
